# extra/clamav/clamd.service
[Unit]
Description=clamav daemon
After=network.target

[Service]
Type=forking
PIDFile=/run/clamav/clamd.pid
ExecStart=/usr/bin/clamd

[Install]
WantedBy=multi-user.target
