# alarm/odroid-c1-remote/odroid_remote.service
[Unit]
Description=ODROID-C1 Remote Control Setup

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/odroid_remote /etc/odroid_remote.conf

[Install]
WantedBy=multi-user.target
