# community/erlang-nox/epmd.socket
[Unit]
Conflicts=epmd.service

[Socket]
ListenStream=4369
Accept=no

[Install]
WantedBy=sockets.target
