# extra/lirc/irexec.service
[Unit]
Description=IR Exec
After=lirc.service
Wants=lirc.service

[Service]
ExecStart=/usr/bin/irexec --daemon
Type=forking

[Install]
WantedBy=multi-user.target
