# extra/lirc/lirc.service
[Unit]
Description=Linux Infrared Remote Control

[Service]
ExecStartPre=/usr/bin/ln -sf /run/lirc/lircd /dev/lircd
ExecStart=/usr/sbin/lircd --pidfile=/run/lirc/lircd.pid
Type=forking
PIDFile=/run/lirc/lircd.pid

[Install]
WantedBy=multi-user.target
