# extra/lirc/lircm.service
[Unit]
Description=LIRC Mouse
After=lirc.service
Wants=lirc.service

[Service]
ExecStart=/usr/sbin/lircmd -u
Type=forking

[Install]
WantedBy=multi-user.target
