# extra/lirc/lircm.service
[Unit]
Description=LIRC Mouse
Requires=lirc.service
After=network.target lirc.service

[Service]
Type=forking
ExecStart=/usr/sbin/lircmd 

[Install]
WantedBy=multi-user.target

