# alarm/alsa-ucm-nyan/90-alsa-ucm-nyan.rules
SUBSYSTEM!="sound", GOTO="sound_end"
ACTION!="change", GOTO="sound_end"
KERNEL!="card*", GOTO="sound_end"

ATTRS{id}=="NVIDIA Tegra Venice2", RUN+="/usr/bin/alsaucm -c NVIDIA\ Tegra\ Venice2 set _verb HiFi"

LABEL="sound_end"
