# extra/slmodem-utils/filelist
.FILELIST
.PKGINFO
etc/
etc/conf.d/
etc/conf.d/slmodem
etc/conf.d/slmodem-alsa
etc/rc.d/
etc/rc.d/slmodem
etc/rc.d/slmodem-alsa
etc/udev/
etc/udev/rules.d/
etc/udev/rules.d/slmodem.rules
usr/
usr/sbin/
usr/sbin/slmodemd
