# community/lxc/lxc-auto.service
[Unit]
Description=Autostart LXC Container
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/lxc-autostart
ExecStop=/usr/bin/lxc-autostart -s
Delegate=true
TasksMax=32768
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
