# aur/umurmur/umurmur.service
[Unit]
Description=Mumble Server
After=network.target

[Service]
PIDFile=/run/umurmur.pid
EnvironmentFile=/etc/conf.d/umurmur
ExecStart=/usr/bin/umurmurd $UMURMUR_ARGS

[Install]
WantedBy=multi-user.target

