# community/opensips/opensips.service
[Unit]
Description=OpenSIPS daemon
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
PIDFile=/run/opensips.pid
ExecStart=/usr/bin/opensips -f /etc/opensips/opensips.cfg -w /var/tmp -P /run/opensips.pid

[Install]
WantedBy=multi-user.target
