# core/nfs-utils/rpc-statd.service
[Unit]
Description=NFSv2/3 Network Status Monitor Daemon
After=rpcbind.service
Requires=rpcbind.service

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/nfs-common.conf
ExecStart=/usr/sbin/rpc.statd $STATD_OPTS

[Install]
WantedBy=multi-user.target
