# core/nfs-utils/rpc-svcgssd.service
[Unit]
Description=RPC GSS-API server-side daemon
After=rpcbind.service
Requires=rpcbind.service

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/nfs-server.conf
ExecStart=/usr/sbin/rpc.svcgssd $SVCGSSD_OPTS

[Install]
WantedBy=multi-user.target
