# aur/lightdm/lightdm-plymouth.service
[Unit]
Description=LightDM Display Manager
After=systemd-user-sessions.service

Conflicts=plymouth-quit.service
After=plymouth-quit.service

[Service]
ExecStart=/usr/sbin/lightdm
StandardOutput=syslog

[Install]
WantedBy=graphical.target
Alias=display-manager.service
