# core/initscripts/adjtime.cron
#!/bin/bash
# Update our hwclock for system drift

. /etc/rc.conf

HWCLOCK_PARAMS="--adjust"
if [ "$HARDWARECLOCK" = "UTC" ]; then
	HWCLOCK_PARAMS="$HWCLOCK_PARAMS --utc"
elif [ "$HARDWARECLOCK" = "localtime" ]; then
	HWCLOCK_PARAMS="$HWCLOCK_PARAMS --localtime"
else
	HWCLOCK_PARAMS=""
fi

if [ -n "$HWCLOCK_PARAMS" ]; then
	/sbin/hwclock $HWCLOCK_PARAMS
fi
