# aur/sickbeard-git/sickbeard-system.service
[Unit]
Description=SickBeard Daemon
After=network.target

[Service]
User=sickbeard
Group=sickbeard
ExecStart=/usr/bin/env python2 /opt/sickbeard/SickBeard.py --quiet --config /opt/sickbeard/config.ini --datadir /opt/sickbeard

[Install]
WantedBy=multi-user.target
