# extra/mysql/mysqld.service
[Unit]
Description=MySQL Server

[Service]
User=mysql
ExecStart=/usr/bin/mysqld --user=mysql
ExecStartPost=/usr/bin/mysqld-post
Restart=always

[Install]
WantedBy=multi-user.target
