# community/syncthing/syncthing@.service
[Unit]
Description=Syncthing service for %i
After=network.target

[Service]
User=%i
Environment=STNORESTART=yes
ExecStart=/usr/bin/syncthing
Restart=always

[Install]
WantedBy=multi-user.target
