# community/syncthing/syncthing.service
[Unit]
Description=Syncthing service for %u
After=network.target

[Service]
Environment=STNORESTART=yes
ExecStart=/usr/bin/syncthing
Restart=on-success

[Install]
WantedBy=default.target
