# community/apcupsd/apcupsd.service
[Unit]
Description=APC UPS Monitor

[Service]
ExecStart=/usr/bin/apcupsd -b

[Install]
WantedBy=multi-user.target
