# extra/courier-mta/courier-webmail-cleancache.cron.hourly
#!/bin/sh

# Cleans the cache of the sqwebmail server
if [ -x /usr/share/sqwebmail/cleancache.pl ]; then
  su -c "/usr/share/sqwebmail/cleancache.pl" bin
fi
