# community/couchdb/couchdb.service
[Unit]
Description=CouchDB Server

[Service]
User=couchdb
Type=simple
Environment="ERL_FLAGS=-couch_ini /usr/lib/couchdb/etc/default.ini /usr/lib/couchdb/etc/datadirs.ini /etc/couchdb/local.ini"
ExecStart=/usr/lib/couchdb/bin/couchdb
ProtectSystem=true
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
