# aur/shinobi-git/shinobi-cron.service
[Unit]
Description=Shinobi cleanup component
After=network-online.target
PartOf=shinobi.target

[Service]
Type=simple
User=shinobi
WorkingDirectory=/usr/share/shinobi
ExecStart=/usr/bin/node /usr/share/shinobi/cron.js

[Install]
WantedBy=shinobi.target
