# aur/shinobi-git/shinobi.target
[Unit]
Description=Shinobi Open Source CCTV and NVR Solution
After=network-online.target
Wants=shinobi-camera.service shinobi-cron.service

[Install]
WantedBy=multi-user.target
