# aur/shinobi-git/shinobi-camera.service
[Unit]
Description=Shinobi server component
After=network-online.target
PartOf=shinobi.target

[Service]
Type=simple
User=shinobi
WorkingDirectory=/usr/share/shinobi
ExecStart=/usr/bin/node /usr/share/shinobi/camera.js

[Install]
WantedBy=shinobi.target
