# core/linux-amlogic-3.14/0hdmi.service
[Unit]
Description=AMlogic HDMI Initialization

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/boot/hdmi.sh

[Install]
WantedBy=multi-user.target
