# aur/open-iscsi/open-iscsi.service
[Unit]
Description=Open iSCSI Daemon
After=network.target

[Service]
Type=forking
ExecStart=/sbin/iscsid

[Install]
WantedBy=multi-user.target
