# alarm/kodi-c2/kodi.service
[Unit]
Description = Kodi for C2
After = remote-fs.target

[Service]
User = kodi
Group = kodi
Type = simple
ExecStart = /usr/bin/kodi
Restart = on-failure

[Install]
WantedBy = multi-user.target
