# aur/tt-rss/tt-rss-updated.service
[Unit]
Description=Tiny tiny rss feed updater daemon
After=network.target

[Service]
ExecStart=/usr/bin/php /usr/share/webapps/tt-rss/update.php -daemon
User=http
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
