# alarm/xbmc-imx/xbmc.service
[Unit]
Description = Starts an instance of XBMC
After = remote-fs.target

[Service]
User = xbmc
Group = xbmc
Type = simple
ExecStart = /usr/bin/xbmc-standalone
Restart = on-failure

[Install]
WantedBy = multi-user.target
