# aur/cherokee/cherokee.service
[Unit]
Description=Cherokee web server
After=network.target

[Service]
ExecStart=/usr/sbin/cherokee

[Install]
WantedBy=multi-user.target
