# aur/gmrender-resurrect-git/gmediarender-user.service
[Unit]
Description=UPnP renderer using gstreamer
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/gmediarender -f "gMedia UpPnP Renderer" -u 49e4f137-a832-481c-9117-bfa96a71aba5

[Install]
WantedBy=default.target
