# aur/gmrender-resurrect-git/gmediarender.service
[Unit]
Description=UPnP renderer using gstreamer
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/gmediarender
ExecStart=/usr/local/bin/gmediarender -f "${friendly}" -u "${uuid}"-d

[Install]
WantedBy=multi-user.target
