# community/gitlab/gitlab-backup.timer
[Unit]
Description=Daily gitlab backup

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=multi-user.target
