# community/gitlab/gitlab-backup.service
[Unit]
Description=GitLab Backup process
Requires=
After=network.target

[Service]
Type=oneshot
User=gitlab
Group=gitlab
SyslogIdentifier=gitlab-backup
WorkingDirectory=<DATADIR>
Environment=RAILS_ENV=production PATH=/opt/ruby/bin:/usr/bin
CapabilityBoundingSet=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
ExecStart=/usr/bin/bundle-2.3 exec rake gitlab:backup:create
