# extra/rspamd/rspamd-dmarc-report.timer
[Unit]
Description=Daily DMARC Reports with rspamd

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

