# extra/rspamd/rspamd-dmarc-report.service
[Unit]
Description=rapid spam DMARC report

[Service]
Type=oneshot
User=rspamd
ExecStart=/usr/bin/rspamadm dmarc_report

