# extra/sane/saned.service
[Unit]
Description=Scanner Service

[Service]
Group=scanner
ExecStart=/usr/bin/saned
StandardInput=socket
StandardError=syslog
