# aur/knockd/knockd.service
[Unit]
Description=Port-Knocking Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/knockd

[Install]
WantedBy=multi-user.target
