# aur/knockd/knockd.service
[Unit]
Description=Port-Knocking Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/knockd -d
Type=forking

[Install]
WantedBy=multi-user.target
