# community/xbmc/xbmc.service
[Unit]
Description = Starts instance of XBMC using xinit
After = systemd-user-sessions.service
Conflicts=getty@tty7.service

[Service]
User = xbmc
Group = xbmc
PAMName=login
Type = simple
TTYPath=/dev/tty7
ExecStart = /usr/bin/xinit /usr/bin/dbus-launch --exit-with-session /usr/bin/xbmc-standalone -- :0 -nolisten tcp vt7
Restart = on-abort

[Install]
WantedBy = multi-user.target
