# community/xbmc/xbmc.service
[Unit]
Description = Starts instance of XBMC using xinit
After = remote-fs.target

[Service]
User = xbmc
Group = xbmc
PAMName=login
Type = simple
ExecStart = /usr/bin/xinit /usr/bin/dbus-launch /usr/bin/xbmc-standalone -l /run/lirc/lircd -- :0 -nolisten tcp
Restart = on-abort

[Install]
WantedBy = multi-user.target
