# community/xbmc/xbmc.service
[Unit]
Description = Starts instance of XBMC using xinit
After = remote-fs.target

[Service]
User = xbmc
Group = xbmc
Type = simple
ExecStart = /usr/bin/xinit /usr/bin/xbmc-standalone -- :0
Restart = always

[Install]
WantedBy = multi-user.target
