# community/squid/squid.service
[Unit]
Description=Web Proxy Cache Server
After=network.target

[Service]
Type=forking
PIDFile=/run/squid.pid
ExecStart=/usr/bin/squid -sYC
ExecStop=/usr/bin/squid -k shutdown
ExecReload=/usr/bin/squid -k reconfigure

[Install]
WantedBy=multi-user.target
