# community/squid/squid-rotate.service
[Unit]
Description=Squid Rotate

[Service]
Type=oneshot
ExecStart=/usr/bin/squid -k rotate
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
