# community/squid/squid-rotate.timer
[Unit]
Description=Weekly rotation of squid log files

[Timer]
OnCalendar=weekly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=timers.target
