# alarm/kodi-rbp-git/kodi.service
[Unit]
Description = Kodi standalone
After = remote-fs.target

[Service]
User = kodi
Group = kodi
Type = simple
ExecStart = /usr/bin/kodi-standalone
Restart = on-failure
KillMode = mixed

[Install]
WantedBy = multi-user.target
