# aur/sabnzbd/sabnzbd.service
[Unit]
Description=SABnzbd binary newsreader
After=network.target

[Service]
Type=simple
ExecStart=/bin/sh -c "python2 /opt/sabnzbd/SABnzbd.py -l0"
User=sabnzbd
Group=sabnzbd

[Install]
WantedBy=multi-user.target
