# alarm/marvell-libvmeta/vmeta-clk@667.service
[Unit]
Description=Set vmeta clock to %I
ConditionFileIsExecutable=/bin/sh
ConditionPathIsReadWrite=/sys/devices/platform/dove_clocks_sysfs.0/vmeta

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo %I000000 > /sys/devices/platform/dove_clocks_sysfs.0/vmeta"
ExecStop=/bin/sh -c "echo 500000000 > /sys/devices/platform/dove_clocks_sysfs.0/vmeta"
RemainAfterExit=yes
