# community/mythtv/mythbackend.service
[Unit]
Description=MythTV Backend
After=network.target mysqld.service

[Service]
Type=simple
Environment=HOME=/var/lib/mythtv
User=mythtv
ExecStart=/usr/bin/mythbackend --logpath /var/log/mythtv
ExecStop=/usr/bin/mythshutdown --setscheduledwakeup
#ExecStop=/usr/bin/mythshutdown --shutdown

[Install]
WantedBy=multi-user.target
