# aur/couchpotato-git/couchpotato.service
[Unit]
Description=An automatic NZB and torrent movie downloader

[Service]
EnvironmentFile=/etc/conf.d/couchpotato_systemd
ExecStart=/usr/bin/python2 /opt/couchpotato/CouchPotato.py $CP_OPTS
GuessMainPID=no
Type=forking
User=couchpotato
Group=couchpotato

[Install]
WantedBy=multi-user.target


