# aur/lirc-user-service/60-lirc.rules
# Make the /dev/lirc* devices accessible for users in the group "lirc"
# using regular group permissions.

KERNEL=="lirc[0-9]*", SUBSYSTEM=="lirc", GROUP="lirc", MODE="0660"

# Grant rw permission to the lirc group for USB devices using acl(5).

ACTION=="add", SUBSYSTEM=="usb", ENV{DEVNAME}=="?*", \
    RUN+="/usr/bin/sh -c '/usr/bin/setfacl -m g:lirc:rw %E{DEVNAME} 2>/dev/null || :'"
