# aur/flexget/flexget.service
[Unit]
Description = FlexGet Daemon
After = network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart = /usr/bin/flexget daemon start
ExecStop = /usr/bin/flexget daemon stop
ExecReload = /usr/bin/flexget daemon reload
KillMode=process

[Install]
WantedBy = default.target
