# extra/samba/smbd@.service
[Unit]
Description=Samba SMB/CIFS server instance

[Service]
ExecStart=/usr/bin/smbd -F
ExecReload=/bin/kill -HUP $MAINPID
StandardInput=socket
