# community/linux-tools/cpupower.service
[Unit]
Description=Apply cpupower configuration

[Service]
Type=oneshot
EnvironmentFile=/etc/default/cpupower
ExecStart=/usr/lib/systemd/scripts/cpupower
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
