# community/linux-tools/usbipd.service
[Unit]
Description=USB/IP server

[Service]
ExecStart=/usr/sbin/usbipd

[Install]
WantedBy=multi-user.target
