# aur/hd-idle/hd-idle.service
[Unit]
Description=Hard drive idling daemon

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/hd-idle
ExecStart=/usr/sbin/hd-idle $HD_IDLE_OPTS

[Install]
WantedBy=multi-user.target
