# aur/input-event-daemon/input-event-daemon.service
# input-event-daemon.service
# Service file for input-event-daemon
# https://github.com/gandro/input-event-daemon

[Unit]
Description=Input Event Daemon

[Service]
Type=forking
ExecStart=/usr/bin/input-event-daemon

[Install]
WantedBy=multi-user.target
