# community/distcc/distccd.service
[Unit]
Description=A distributed C/C++ compiler
Documentation=man:distccd(1)

[Service]
EnvironmentFile=/etc/conf.d/distccd
ExecStart=/usr/bin/distccd --daemon --no-detach $DISTCC_ARGS

[Install]
WantedBy=multi-user.target
