# extra/netperf/netserver.service
[Unit]
Description=Netperf Benchmark Server
After=network.target

[Service]
ExecStart=/usr/bin/netserver -D

[Install]
WantedBy=multi-user.target
