# core/klibc-udev/50-udev-default.rules
# firmware class requests
SUBSYSTEM=="firmware", ACTION=="add", RUN+="firmware.sh"

# do not delete static device nodes
ACTION=="remove", NAME=="?*", TEST=="/lib/udev/devices/$name", OPTIONS+="ignore_remove"
