# extra/ceph/ceph-osd@.service
[Unit]
Description=Ceph OSD %i
After=network.target

[Service]
ExecStart=/usr/bin/ceph-osd -f -i %i

[Install]
WantedBy=multi-user.target
