# extra/ceph/ceph-mds@.service
[Unit]
Description=Ceph MDS %i
After=network.target

[Service]
ExecStart=/usr/bin/ceph-mds -f -i %i

[Install]
WantedBy=multi-user.target
