# community/uwsgi/uwsgi_at.service
[Unit]
Description=uWSGI service unit
After=syslog.target

[Service]
ExecStart=/usr/bin/uwsgi --ini /etc/uwsgi/%I.ini
ExecReload=/usr/bin/uwsgi --ini /etc/uwsgi/%I.ini --reload $MAINPID
ExecStop=-/usr/bin/uwsgi --ini /etc/uwsgi/%I.ini --stop $MAINPID
Restart=always
Type=notify
StandardError=syslog
NotifyAccess=main
KillSignal=SIGQUIT

[Install]
WantedBy=multi-user.target
