# community/uwsgi/uwsgi@.service
[Unit]
Description=uWSGI service unit
After=syslog.target

[Service]
ExecStart=/usr/bin/uwsgi --ini /etc/uwsgi/%I.ini
ExecReload=/usr/bin/uwsgi --ini /etc/uwsg/%I.ini --reload
ExecStop=/usr/bin/uwsgi --ini /etc/uwsg/%I.ini --stop
Restart=always
Type=notify
StandardError=syslog
NotifyAccess=main
KillSignal=SIGQUIT

[Install]
WantedBy=multi-user.target
