# extra/hylafax/hylafax.cron.daily
#!/bin/sh

/usr/sbin/faxqclean
/usr/sbin/faxcron -rcv 30
