# alarm/fake-hwclock/fake-hwclock.service
[Unit]
Description=Restore system time on boot and save it on shutdown
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/fake-hwclock.sh load
ExecStop=/usr/lib/systemd/scripts/fake-hwclock.sh save
RemainAfterExit=true

[Install]
WantedBy=sysinit.target
