# alarm/fake-hwclock/fake-hwclock-save.service
[Unit]
Description=Periodically saves system time to file
After=fake-hwclock.service
Requires=fake-hwclock.service

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/fake-hwclock.sh save
StandardOutput=null
