# alarm/fake-hwclock/fake-hwclock.service
[Unit]
Description=Sets time to what it was at last shutdown.
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/savesettime.sh set
RemainAfterExit=true
ExecStop=/usr/lib/systemd/scripts/savesettime.sh save

[Install]
WantedBy=sysinit.target

